# visitors_addon/files/php83/etc/services.d/php-fpm83/finish
#!/usr/bin/with-contenv sh
# Always exit cleanly; prevents "Exec format error" complaints
exit 0
