# visitors_addon/files/php83/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh
set -e
mkdir -p /run/nginx /run/php
chown -R nginx:nginx /run/nginx /run/php
exec nginx -g 'daemon off;'
